# rust/workplacesim/deploy/image-overlay/etc/inittab
# /etc/inittab — Alpine Linux default, with tty1 getty disabled so the
# workplacesim renderer can own the VT and draw to /dev/fb0 without a
# login prompt flickering over the top.
#
# This is a static copy of Alpine's stock rpi inittab (as shipped in the
# alpine-rpi tarball). If Alpine changes the default in a future release
# this overlay file will overwrite their new default — keep an eye on the
# upstream when bumping the Alpine version in mkimage.sh.

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# Set up a couple of getty's
#tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4
tty5::respawn:/sbin/getty 38400 tty5
tty6::respawn:/sbin/getty 38400 tty6

# Put a getty on the serial port
ttyAMA0::respawn:/sbin/getty -L ttyAMA0 115200 vt100

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/sbin/openrc shutdown
